# kdlp-dog.service
[Unit]
Description=KDLP ORBIT externally-hosted dog-themed service degradation monitor 

[Service]
ExecStart=/var/orbit/dog.py/dogdaemon.py
